# config/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:iventoy]
directory=/iventoy
command=/bin/sh -c "env IVENTOY_API_ALL=1 IVENTOY_AUTO_RUN=1 /iventoy/lib/iventoy"
